# contrib/nfpm/preremove.sh
#!/bin/sh
systemctl --global disable z13gui.service || true
