# contrib/z13gui.service
[Unit]
Description=z13gui - GTK4 overlay drawer for z13ctl
Documentation=https://github.com/dahui/z13gui
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/z13gui
Environment=GTK_A11Y=none
EnvironmentFile=-%t/gamescope-environment
Restart=on-failure
RestartSec=3
ExecStopPost=/bin/sh -c 'kill -CONT $(cat /run/user/%U/z13gui-frozen-pid 2>/dev/null) 2>/dev/null; rm -f /run/user/%U/z13gui-frozen-pid'

[Install]
WantedBy=graphical-session.target
