# contrib/aur/z13gui-bin.install
post_install() {
    udevadm control --reload-rules
    udevadm trigger
    systemctl --global enable z13gui.service
}

post_upgrade() {
    udevadm control --reload-rules
    udevadm trigger
}
