# contrib/aur/z13gui-bin.install
post_install() {
    udevadm control --reload-rules
    udevadm trigger
    setcap cap_bpf,cap_perfmon+ep /usr/bin/z13gui 2>/dev/null || true
    systemctl --global enable z13ctl.socket z13ctl.service
    systemctl --global enable z13gui.service
}

post_upgrade() {
    udevadm control --reload-rules
    udevadm trigger
    setcap cap_bpf,cap_perfmon+ep /usr/bin/z13gui 2>/dev/null || true
}
